# fans/crontab
0 7 * * * i2cset -y 1 0x01a 0xfe
0 20 * * * i2cset -y 1 0x01a 0x00
